# app/.scratch/etc/s6-overlay/s6-rc.d/app/run
#!/command/execlineb -P
with-contenv
cd /opt/app
/opt/app/entrypoint.sh
